# config/crontab
0 0 * * 0 /usr/local/bin/update_and_reboot.sh
